# rocky/playbooks/roles/k8s_common/files/system/kubernetes.slice
[Unit]
Description=kubernetes slice
Wants=-.slice

[Slice]

[Unit]
DefaultDependencies=no
