# centos/playbooks/roles/k8s_common/files/system/kubernetes.slice
# http://www.jinbuguo.com/systemd/systemd.slice.html
# http://www.jinbuguo.com/systemd/systemd.unit.html
# http://www.jinbuguo.com/systemd/systemd.service.html
[Unit]
Description=kubernetes slice
DefaultDependencies=no
Wants=-.slice

# http://www.jinbuguo.com/systemd/systemd.resource-control.html
[Slice]
